# 2_Packages/mudskipper_mod/etc/profile.d/mudskipper.sh
#	This Script will run on Login

#########################
#Mudskipper Extensions  #
#########################

/bin/bash /etc/mudskipper/profile.sh



/etc/mudskipper/Themes/script/theme

